# files/etc/uci-defaults/999-custom-settings
## /etc/uci-defaults/999-custom-settings
## removed after execution with exit 0

sed -E -e '/routing|telephony|video/s/^/#/' \
    -e 's|mirrors.vsean.net/openwrt|mirrors.pku.edu.cn/immortalwrt|' \
    -i /etc/apk/repositories.d/distfeeds.list

exit 0
